# init-scripts/gentoo/schedtoold
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2


depend() {
	use logger
	}

start() {
	ebegin "Starting schedtoold - a schedtool daemon"
	start-stop-daemon --start --quiet --background - \
	--exec /usr/bin/schedtoold -- -p /var/run/schedtoold.pid
	eend $?
}

stop() {
	ebegin "Stopping schedtoold - a schedtool daemon"
	start-stop-daemon --stop --quiet --pidfile /var/run/schedtoold.pid
	eend $?
}
